# infra/systemd/ollama-dev1.service
# Ollama (dev1)
# Place into /etc/systemd/system/ollama-dev1.service
# Then: systemctl daemon-reload && systemctl enable --now ollama-dev1
[Unit]
Description=Ollama LLM Server (dev1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/dev1/ollama
Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=OLLAMA_MODELS=/mnt/ai_storage/models/ollama
Environment=OLLAMA_ORIGINS=*
ExecStart=/usr/local/bin/ollama serve
User=ollama
Group=ollama
Restart=always
RestartSec=3
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin"

[Install]
WantedBy=multi-user.target
